# deploy/supervisor.conf
[program:videoforge-api]
command=/var/www/videoforge/venv/bin/uvicorn server.app:app --host 127.0.0.1 --port 8000 --workers 2
directory=/var/www/videoforge
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/videoforge/error.log
stdout_logfile=/var/log/videoforge/access.log
environment=PATH="/var/www/videoforge/venv/bin",NODE_ENV="production",APP_ENV="production"
